# ops/com.argus.pages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DORMANT until go-live. Copy to ~/Library/LaunchAgents/ and `launchctl load`
     it ONLY after GitHub Pages is enabled — see ops/GO_LIVE.md. Redeploys the
     public terminal every 10 min so the live site tracks the book. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.argus.pages</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/agent/argus_v2/ops/deploy_pages.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/argus_pages.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/argus_pages.log</string>
</dict>
</plist>
